# systemd/animaksm.service
[Unit]
Description=Animaksm KSM Governor and Memory Scanner
Documentation=https://github.com/animaios/animaksm
After=local-fs.target
ConditionPathExists=/sys/kernel/mm/ksm

[Service]
Type=simple
ExecStart=/usr/local/bin/animaksm run --config /etc/animaksm.toml
ExecStopPost=-/usr/local/bin/animaksm restore-ksm --config /etc/animaksm.toml
Restart=on-failure
RestartSec=5

# Security: limit what the process can do outside KSM tuning.
NoNewPrivileges=yes
AmbientCapabilities=CAP_SYS_PTRACE CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_SYS_PTRACE CAP_DAC_READ_SEARCH
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
# /sys/kernel/mm/ksm must be read-write for governor + scanner
ReadWritePaths=/var/lib/animaksm /sys/kernel/mm/ksm
ReadOnlyPaths=/proc/pressure /proc /etc/animaksm.toml
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=no
RemoveIPC=yes

# Resource limits
MemoryMax=32M
Nice=19

[Install]
WantedBy=multi-user.target
